# debian/postinst
#!/bin/sh

# Exit on error
set -e

systemctl start indoor_sensors.service
systemctl enable indoor_sensors.service